# templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User={{ user }}
ExecStart={{ prom_dir }}/node_exporter

[Install]
WantedBy=default.target

